# site-cookbooks/round-three/templates/default/faye-initd.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          faye
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
# Short-Description: faye initscript
# Description:       faye
### END INIT INFO

# Original author: Forrest Robertson

# Do NOT "set -e"

case "$1" in
  start)
  cd <%= @round_three_dir %>
  <%= @bundler %> exec thin start --ssl --ssl-verify --ssl-key-file /etc/ssl/transit-demo.admin.umass.edu.key
--ssl-cert-file /etc/ssl/transit-demo.admin.umass.edu.cert -C <%= @config %>
	;;
  stop)
  cd <%= @round_three_dir %>
  <%= @bundler %> exec thin stop -C <%= @config %>
	;;
  restart)
  cd <%= @round_three_dir %>
  <%= @bundler %> exec thin restart -C <%= @config %>
	;;
  *)
	echo "Usage: $SCRIPT_NAME {start|stop|restart}" >&2
	exit 3
	;;
esac
:
